# app/busybox/files/etc/inittab
::sysinit:/etc/init.d/rcS
ttyS0::respawn:-/bin/sh
tty1::respawn:-/bin/sh
::ctrlaltdel:/bin/umount -a -
